# puppet/manifests/files/etc/profile.d/git-prompt.sh
export LC_ALL=en_US.UTF-8
export LANG=en_US.UTF-8
. /usr/local/share/bash-git-prompt/gitprompt.sh 
